# launch.service
[Unit]
Description=Resin.io User Application

[Service]
WorkingDirectory=/usr/bin/
EnvironmentFile=/etc/docker.env
ExecStart=/etc/resinApp.cmd
After=systemd-udev-trigger.service

[Install]
WantedBy=basic.target